# monitor-test.service
[Unit]
Description=Мониторинг процесса test

[Service]
ExecStart=/usr/local/bin/monitor_test.sh